# target/etc/systemd/system/bespoke-toggle-interface.service
I have the follwing system service. Can I be sure that toffle_interface.sh will be run as root bu default?

[Unit]
Description=Prepare Bespoke to run on wlan0 or wlan1
After=network-online.target
Wants=network-online.target

[Service]
User=root
Type=oneshot
ExecStart=/home/pi/bespoke/target/bin/toggle_interface.sh

[Install]
WantedBy=multi-user.target
