# target/etc/systemd/system/hostapd-wlan1.service
[Unit]
Description=Access point and authentication server for Wi-Fi and Ethernet (wlan1)
Documentation=man:hostapd(8)
After=network.target
ConditionFileNotEmpty=/etc/hostapd/hostapd-wlan1.conf

[Service]
Type=forking
PIDFile=/run/hostapd-wlan1.pid
Restart=on-failure
RestartSec=2
Environment=DAEMON_CONF=/etc/hostapd/hostapd-wlan1.conf
EnvironmentFile=-/etc/default/hostapd
ExecStart=/usr/sbin/hostapd -B -P /run/hostapd-wlan1.pid $DAEMON_OPTS ${DAEMON_CONF}

[Install]
WantedBy=multi-user.target
