# build/config/bespoke.service
[Unit]
Description=Bespoke BBS
After=network.target

[Service]
ExecStart=/bin/bash -c '/home/pi/bespoke/bin/bespoke --release'
WorkingDirectory=/home/pi/bespoke
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
